# overlay/etc/supervisord.d/sniproxy.conf
[program:sniproxy]
command=sniproxy -f
autostart=false
