# overlay/etc/supervisord.d/sniproxy.conf
[program:sniproxy]
command=sniproxy -f
autostart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0