# overlay/app/dnslocal.conf
[program:dnslocal]
command=dnsmasq -k -C /etc/dnsmasq_local.conf
autostart=false
