# overlay/app/init.conf
[program:init]
command=python /app/init.py
autorestart=false